# systemd/jack-site-caddy.service
[Unit]
Description=jack.kelliher.info — Caddy static file server
After=network.target

[Service]
Type=notify
ExecStart=/usr/bin/env caddy run --config /etc/jack-site/Caddyfile
ExecReload=/usr/bin/env caddy reload --config /etc/jack-site/Caddyfile
Restart=on-failure
RestartSec=5

# Sandboxing
DynamicUser=yes
ReadOnlyPaths=/
ReadWritePaths=/var/lib/caddy
ProtectHome=yes
PrivateTmp=yes
NoNewPrivileges=yes
ProtectSystem=strict

# Caddy needs to bind to port 8780
AmbientCapabilities=

[Install]
WantedBy=multi-user.target
