# systemd/jack-site-cloudflared.service
[Unit]
Description=jack.kelliher.info — Cloudflare Tunnel
After=network-online.target jack-site-caddy.service
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/jack-site/tunnel-env
ExecStart=/usr/bin/env cloudflared --no-autoupdate tunnel run --token ${TUNNEL_TOKEN}
Restart=on-failure
RestartSec=10

# Sandboxing
DynamicUser=yes
ProtectHome=yes
PrivateTmp=yes
NoNewPrivileges=yes
ProtectSystem=strict

[Install]
WantedBy=multi-user.target
